# angela_core/backup/com.david.angela.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.david.angela.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3</string>
        <string>/Users/davidsamanyaporn/PycharmProjects/AngelaAI/angela_core/backup/run_backup.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/davidsamanyaporn/PycharmProjects/AngelaAI</string>

    <!-- Run daily at 3:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Also run at load (for testing) -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/davidsamanyaporn/PycharmProjects/AngelaAI/logs/backup_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/davidsamanyaporn/PycharmProjects/AngelaAI/logs/backup_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PYTHONPATH</key>
        <string>/Users/davidsamanyaporn/PycharmProjects/AngelaAI</string>
    </dict>
</dict>
</plist>
